# systemd/restaurant-enhanced-api.service
[Unit]
Description=Restaurant Enhanced API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspace/projects
Environment="PATH=/workspace/projects/venv/bin"
ExecStart=/workspace/projects/venv/bin/python -m uvicorn src.api.restaurant_enhanced_api:app --host 0.0.0.0 --port 8007
Restart=always
RestartSec=10
StandardOutput=append:/workspace/projects/logs/enhanced_api.log
StandardError=append:/workspace/projects/logs/enhanced_api.error.log

[Install]
WantedBy=multi-user.target
