# systemd/restaurant-api.service
[Unit]
Description=Restaurant API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspace/projects
Environment="PATH=/workspace/projects/venv/bin"
ExecStart=/workspace/projects/venv/bin/python -m uvicorn src.api.restaurant_api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/workspace/projects/logs/api.log
StandardError=append:/workspace/projects/logs/api.error.log

[Install]
WantedBy=multi-user.target
